# systemd/DasConcatDaily.service
[Unit]
Description=Run concatenation of FEBUS data

[Service]
ExecStart=/bin/bash PROJECT_PATH/concat.sh

[Install]
WantedBy=multi-user.target