# systemd/DasConcatDaily.timer
[Unit]
Description=Run concatenation of FEBUS data daily at 3am.

[Timer]
OnCalendar=*-*-* 3:00:00
Persistent=true

[Install]
WantedBy=timers.target